# debian/postinst
#!/bin/bash
systemctl enable wolkconnect-task
systemctl start wolkconnect-task